# bin/php83/crontab
# Example cron job: Run a script every minute
* * * * * root echo "Cron job running" >> /var/log/cron.log 2>&1

# Add more cron jobs as needed